# pkg/DEBIAN/prerm
#!/bin/bash

rm /etc/nginx/sites-enabled/ble2json.conf

systemctl stop ble2json.service

systemctl disable ble2json.service
